# systemd/sample-app.service
[Unit]
Description=Sample Instrumented App — emits traces to Tempo via OpenTelemetry
After=network-online.target otel-collector.service
Wants=network-online.target

[Service]
User=sampleapp
Group=sampleapp
Type=simple
WorkingDirectory=/opt/sample-app
Environment="OTEL_EXPORTER_OTLP_ENDPOINT=http://localhost:4317"
Environment="OTEL_SERVICE_NAME=sample-app"
Environment="OTEL_RESOURCE_ATTRIBUTES=service.version=1.0.0,deployment.environment=production"
ExecStart=/opt/sample-app/venv/bin/uvicorn main:app \
  --host 0.0.0.0 \
  --port 8000
Restart=always
RestartSec=5
SyslogIdentifier=sample-app

[Install]
WantedBy=multi-user.target
